# rootfs/etc/supervisor.d/default.ini
[program:httpd]
command=/usr/local/bin/httpd-foreground 
numprocs=1
autostart=true
autorestart=true

[program:mariadb]
command=/usr/bin/mysqld 
numprocs=1
autostart=true
autorestart=true
