# rootfs/etc/supervisor.d/mariadb.ini
[program:httpd]
user=root
command=/usr/local/bin/httpd-foreground 
numprocs=1
autostart=true
autorestart=true

[program:mariadb]
user=mysql
command=/usr/bin/mysqld 
numprocs=1
autostart=true
autorestart=true
